# srcs/phpmyadmin/supervisor.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx
killasgroup=true
stopasgroup=true
redirect_stderr=true

[program:php-pfm7]
command=php-fpm7

[program:telegraf]
command= /telegraf/telegraf