# srcs/ftps/supervisor.conf
[supervisord]
nodaemon=true

[program:ftps]
command=vsftpd /etc/vsftpd.conf

[program:telegraf]
command= /telegraf/telegraf